# templates/default/rtorrent.user.upstart.erb
description "<%= @user.name %> rtorrent"

expect fork
kill timeout 30

start on (local-filesystems and net-device-up and runlevel [2345])
stop on runlevel [016]

script
  su <%= @user.name %> -c "screen -dms rtorrent /usr/local/bin/rtorrent"
end script

post-stop script
  [ -f <%= @user.home %>/.session/rtorrent.lock ] && kill $(cat <%= @user.home %>/.session/rtorrent.lock | awk -F+ '{ print $2 }')
  [ -f <%= @user.home %>/.session/rtorrent.socket ] && rm <%= @user.home %>/.session/rtorrent.socket
end script
