# templates/default/rtorrent.user.upstart.erb
description "<%= @user.name %> rtorrent"

start on ((local-filesystems and net-device-up and runlevel [2345]) or starting rtorrent)
stop on (runlevel [016] or stopping rtorrent)

pre-start script
  rm -f <%= @user.home %>/.session/rtorrent.socket
  su <%= @user.name %> -c "screen -dms rtorrent /usr/local/bin/rtorrent"
end script

post-stop script
  [ -f <%= @user.home %>/.session/rtorrent.lock ] && kill $(cat <%= @user.home %>/.session/rtorrent.lock | awk -F+ '{ print $2 }')
  rm -f <%= @user.home %>/.session/rtorrent.lock
  rm -f <%= @user.home %>/.session/rtorrent.socket
end script
